# curator4-895180668.conf
[program:curator4-895180668]
command=python3 /home/curator4/curator4-895180668/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator4/curator4-895180668/test.err.log
stdout_logfile=/home/curator4/curator4-895180668/test.out.log